# csye6225-cloudwatch.service
[Unit]
Description=Unit for starting cloudwatch agent
After=network.target

[Service]
ExecStart=/home/ubuntu/scripts/csye6225-cloudwatch.sh
User=ubuntu
Group=ubuntu
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
